# cron.d/letsencrypt
0 */12 * * * certbot renew --nginx --noninteractive --renew-hook "nginx -s reload && /app/comsave/sync_push_certs.sh" >/dev/null 2>&1
