# supervisor/s3_mirror.conf
[program:s3_mirror_certs_push]
command=/usr/local/bin/mc mirror --watch --overwrite --remove /etc/letsencrypt wasabi/nginx-proxy-letsencrypt
autostart=true
autorestart=true
